# octodetect/rules/unplugged.rules
ACTION=="remove", KERNEL=="tty*", RUN+="/bin/sh -c '/app/octoprint_usb_unplugged'"
